# spammenot/xinetd.d/spammenot_pop3
#
# spammenot - pop3 mail daemon
#
service pop3
{
    disable	= no
    socket_type	= stream
    protocol	= tcp
    wait	= no
    user	= root
    server	= /usr/bin/perl 
    server_args	= /etc/spammenot/pop3.pl 
    flags	= IPv4
}
